# freshrss/root/etc/s6-overlay/s6-rc.d/init-prepare/run
#!/command/with-contenv bash

/usr/local/bin/frankenphp php-cli /opt/freshrss/cli/prepare.php > /dev/null