# freshrss/root/etc/s6-overlay/s6-rc.d/svc-supercronic/run
#!/command/with-contenv bash

/usr/local/bin/supercronic /etc/crontab