# freshrss/root/etc/s6-overlay/s6-rc.d/svc-freshrss/run
#!/command/with-contenv bash

/usr/local/bin/frankenphp run --config /opt/freshrss/Caddyfile